# ctf/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:ttyd]
command=/usr/bin/ttyd --writable -p 7681 sudo -iu ctf
autostart=true
autorestart=true
